# noip-renewer-s6/rootfs/etc/s6-overlay/s6-rc.d/noip-renewer/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================

workdir /app

CMD ["python3", "renew.py"]
